# deploy/gateway/pp2050-gw.service
[Unit]
Description=Plan Peru 2050 — secure AI gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/pp2050-gw/gateway.py
Restart=on-failure
RestartSec=3
# La key de OpenRouter vive SOLO aquí (no en el Caddyfile ni en el cliente).
Environment=OPENROUTER_KEY=REEMPLAZAR_KEY
Environment=PP2050_DATA=/var/www/plan-peru-2050/data
Environment=PP2050_GW_PORT=3501
# Modelo principal Gemini 2.5 Flash (barato) + fallback al llama (si Gemini falla).
Environment=PP2050_MODELS=google/gemini-2.5-flash,meta-llama/llama-3.3-70b-instruct
# Alternativa: Gemini NATIVO (aún más barato, como liti.app). Descomenta y pon la key de Google AI Studio:
# Environment=GEMINI_API_KEY=AIza...
# Environment=GEMINI_MODEL=gemini-flash-latest
# Hardening básico
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pp2050-gw

[Install]
WantedBy=multi-user.target
